# lock.service
[Unit]
Description=Lock

[Service]
User=leonardo
Environment=DISPLAY=:0
ExecStart=/usr/local/bin/slock

[Install]
WantedBy=suspend.target
